# etc-s6-overlay/s6-rc.d/otbr-agent/run
#!/command/with-contenv /bin/bash

set -eu

OTBR_REST_LISTEN_ADRESS=${OTBR_REST_LISTEN_ADRESS:-0.0.0.0}
OTBR_REST_LISTEN_PORT=${OTBR_REST_LISTEN_PORT:-8081}
OTBR_THREAD_IF=${OTBR_THREAD_IF:-wpan0}
OTBR_BACKBONE_IF=${OTBR_BACKBONE_IF:-eth0}
OTBR_LOG_LEVEL_INT=${OTBR_LOG_LEVEL_INT:-6}

verifyenvvars="OTBR_REST_LISTEN_ADRESS OTBR_REST_LISTEN_PORT OTBR_THREAD_IF OTBR_BACKBONE_IF"

echo "Startup configuration for otbr-agent:"
for var in $verifyenvvars; do
    value="${!var:-}"

    if [[ -z "$value" ]]; then
        echo "ERROR: Environment variable '$var' is not set or empty"
        exit 1
    fi

    echo "INFO: $var=$value"
done

SOCAT_STARTUP_GRACE_PERIOD=${SOCAT_STARTUP_GRACE_PERIOD:-2}

echo "Sleeping for ${SOCAT_STARTUP_GRACE_PERIOD}s to give socat time to come up ..."
sleep "${SOCAT_STARTUP_GRACE_PERIOD}"
echo "alright, let's move on"

exec 
    "/usr/sbin/otbr-agent" -I ${OTBR_THREAD_IF} -B "${OTBR_BACKBONE_IF}" \
        --rest-listen-address "${OTBR_REST_LISTEN_ADRESS}" \
        --rest-listen-port "${OTBR_REST_LISTEN_PORT}" \
        -d${OTBR_LOG_LEVEL_INT} -v -s \
        "spinel+hdlc+uart://${RCP_TTY:-/tmp/ttyOTBR}?uart-baudrate=${RCP_BAUDRATE:-460800}${RCP_ADDITIONAL_ARGS:-}" \
        "trel://${OTBR_BACKBONE_IF}"
