# etc/s6-overlay/s6-rc.d/dbus/run
#!/bin/bash -eu

echo "Starting dbus-daemon ..."
mkdir -p /var/run/dbus

exec /usr/bin/dbus-daemon \
    --system \
    --nofork \
    --nopidfile \
    --syslog
