# etc/s6-overlay/s6-rc.d/otbr-agent/finish
#!/command/with-contenv /bin/bash
exit 0
