# etc/s6-overlay/s6-rc.d/otbr-agent/finish
#!/command/with-contenv bash

# s6 passes the service exit info as $1 (exit code, or 256 if killed by signal)
# and $2 (signal number when $1 == 256). Translate that to a conventional
# 128+signal exit code, then bring the whole container down on anything other
# than a clean exit so docker's restart policy can take over. Without this,
# s6-supervise would silently restart otbr-agent forever even when the radio
# link or configuration is broken.

sig=0
if test "$1" -eq 256 ; then
	sig="$2"
	e=$((128 + $2))
	echo "otbr-agent was killed by signal ${sig}."
else
	e="$1"
	echo "otbr-agent exited with code ${e}."
fi

# SIGTERM (15) and SIGINT (2) are the orderly-shutdown signals: `docker stop`
# sends SIGTERM, and `s6-svc -d` / Ctrl-C use them too. Treat these as an
# intentional stop -- let s6's down state govern restart and let the container
# exit cleanly -- instead of misreporting a deliberate shutdown as a crash and
# recording a 143/130 exit code.
if test "$sig" -eq 15 || test "$sig" -eq 2 ; then
	exit 0
fi

# Exit code 5 is the upstream convention for "shut down on purpose, don't
# treat as fatal" -- mirror it so an operator can stop the agent cleanly.
if test "$e" -ne 0 && test "$e" -ne 5; then
	echo "$e" > /run/s6-linux-init-container-results/exitcode
	/run/s6/basedir/bin/halt
	exit 125
fi
