# etc/s6-overlay/s6-rc.d/otbr-agent/run
#!/command/with-contenv /bin/bash

set -eu

# otbr-agent logs to stderr (-v); fold it into stdout so the otbr-agent-log
# logger service (see producer-for) captures it as whole, atomic lines.
exec 2>&1

. /usr/local/lib/otbr/booleans.sh

if otbr_is_true "${RCP_USE_TCP:-}" ; then
	echo -n "Sleeping for ${SOCAT_STARTUP_GRACE_PERIOD}s to give socat time to come up ... "
	sleep "${SOCAT_STARTUP_GRACE_PERIOD}"
	echo "Proceeding."
fi

OTBR_REST_LISTEN_ADDRESS="${OTBR_REST_LISTEN_ADDRESS:-${OTBR_REST_LISTEN_ADRESS:-0.0.0.0}}"

echo "Starting otbr-agent ..."

exec \
    "/usr/sbin/otbr-agent" -I "${OTBR_THREAD_IF}" -B "${OTBR_BACKBONE_IF}" \
        --vendor-name "${OTBR_VENDOR_NAME}" \
        --model-name "${OTBR_MODEL_NAME}" \
        --rest-listen-address "${OTBR_REST_LISTEN_ADDRESS}" \
        --rest-listen-port "${OTBR_REST_LISTEN_PORT}" \
        "-d${OTBR_LOG_LEVEL_INT}" -v -s \
        "spinel+hdlc+uart://${RCP_TTY:-/tmp/ttyOTBR}?uart-baudrate=${RCP_BAUDRATE:-460800}${OTBR_RCP_ADDITIONAL_ARGS:-}" \
        "trel://${OTBR_BACKBONE_IF}"
