# etc/s6-overlay/s6-rc.d/otbr-web-log/run
#!/command/with-contenv sh
# Consume otbr-web's piped output and re-emit each line atomically to the
# container stdout (docker logs), prefixed with an ISO 8601 UTC timestamp.
exec /command/s6-log -- T 1
