# etc-s6-overlay/s6-rc.d/otbr-dataset-backup/run
#!/command/with-contenv /bin/bash

while :; do
	/usr/local/bin/otbr-dataset.sh backup
	sleep ${OTBR_BACKUP_INTERVAL:-600}
done
