# etc/s6-overlay/s6-rc.d/otbr-dataset-backup/run
#!/command/with-contenv /bin/bash

while :; do
	sleep ${OTBR_BACKUP_INTERVAL}
	echo "Creating dataset backup ..."
	/usr/local/bin/otbr-dataset.sh backup
done
