# etc/s6-overlay/s6-rc.d/otbr-web/run
#!/command/with-contenv /bin/bash

set -eu

# otbr-web logs to stderr; fold it into stdout so the otbr-web-log logger
# service (see producer-for) captures it as whole, atomic lines.
exec 2>&1

. /usr/local/lib/otbr/booleans.sh

if ! otbr_is_true "${OTBR_WEB_ENABLE:-}" ; then
	echo "otbr-web disabled."
	exec sleep infinity
fi

# Tell otbr-web which port the REST API listens on (-P) so its frontend
# queries the right one. Replaces the old hack of sed-patching app.js.
exec /usr/sbin/otbr-web \
	-I "${OTBR_THREAD_IF}" \
	-d${OTBR_LOG_LEVEL_INT} \
	-s \
	-p "${OTBR_WEB_PORT}" \
	-a "${OTBR_WEB_LISTEN_ADDRESS}" \
	-P "${OTBR_REST_LISTEN_PORT}"
