# etc/s6-overlay/s6-rc.d/otbr-web/run
#!/command/with-contenv /bin/bash

set -eu

# otbr-web logs to stderr; fold it into stdout so the otbr-web-log logger
# service (see producer-for) captures it as whole, atomic lines.
exec 2>&1

if [[ "${OTBR_WEB_ENABLE}" -eq "0" ]] ; then
	echo "otbr-web disabled."
	exec sleep infinity
fi

if [[ "${OTBR_WEB_PATCH_REST_PORT}" -ge "1" ]] ; then
	sed -i "s/'8081'/'${OTBR_REST_LISTEN_PORT}'/" /usr/share/otbr-web/frontend/res/js/app.js
fi

exec /usr/sbin/otbr-web \
	-I "${OTBR_THREAD_IF}" \
	-d${OTBR_LOG_LEVEL_INT} \
	-s \
	-p "${OTBR_WEB_PORT}" \
	-a "${OTBR_WEB_LISTEN_ADDRESS}"
