# etc/s6-overlay/s6-rc.d/socat/run
#!/command/with-contenv /bin/bash

set -eu

# socat's diagnostics go to stderr; fold them into stdout so the socat-log
# logger service (see producer-for) captures them as whole, atomic lines
# instead of letting them collide with otbr-agent's output in docker logs.
exec 2>&1

if [[ "${RCP_USE_TCP}" -eq "0" ]] ; then
	echo "RCP over tcp (socat) disabled."
	exec sleep infinity
fi

if [[ -z "${RCP_HOST:-}" ]]; then
    echo "ERROR: RCP_HOST environment variable is not set"
    exit 1
fi

echo "Starting socat ..."

exec /usr/bin/socat -d \
  pty,link=${RCP_TTY}${SOCAT_SOURCE_PARAMETERS} \
  tcp:${RCP_HOST}:${RCP_PORT}${SOCAT_DESTINATION_PARAMETERS}
