# etc-s6-overlay/s6-rc.d/otbr-dataset-restore/up
#!/command/with-contenv /bin/bash
exec /usr/local/bin/otbr-dataset.sh restore
